# openwrt_builder/rootfs/etc/s6-overlay/s6-rc.d/openwrt-builder-runner/run
#!/command/with-contenv sh
set -eu

set -a
. /etc/openwrt-builder/paths.env
set +a

exec python3 -m openwrt_builder.runner.main
