# openwrt_builder/rootfs/etc/s6-overlay/s6-rc.d/openwrt-builder-runner/run
#!/command/with-contenv sh
set -eu

exec sh /usr/local/bin/openwrt-builder-launcher.sh runner
