# openwrt_builder/rootfs/etc/services.d/openwrt_builder/run
#!/command/with-contenv sh

mkdir -p /data/profiles /data/artifacts /cache

set -euo pipefail
python3 /usr/local/bin/web_api.py &

printf "\033[32mINFO\033[0m OpenWRT ImageBuilder started\r"