# openwrt_builder/rootfs/etc/services.d/openwrt_builder/run
#!/command/with-contenv sh
printf "\033[32mINFO\033[0m OpenWRT ImageBuilder started"

set -euo pipefail
python3 /usr/local/bin/web_api.py &

sleep infinity