# openwrt_builder/rootfs/etc/services.d/openwrt_builder/run
#!/command/with-contenv sh
echo "OpenWRT ImageBuilder started"

set -euo pipefail
python3 /usr/local/bin/health_server.py &

sleep infinity