# openwrt_builder/rootfs/etc/s6-overlay/s6-rc.d/openwrt-builder/run
#!/command/with-contenv sh
printf "\033[32mINFO\033[0m Starting OpenWRT ImageBuilder\n"

set -eu
exec sh /usr/local/bin/openwrt-builder-launcher.sh api
