# openwrt_builder/rootfs/etc/s6-overlay/s6-rc.d/openwrt-builder-init/up
#!/command/with-contenv sh
printf "\033[32mINFO\033[0m INIT step.\n" >&2

# export env paths
set -a
. /etc/openwrt-builder/paths.env
set +a
printf "\033[32mINFO\033[0m paths.env step.\n" >&2

#env | while IFS='=' read -r k v; do
#  case "$k" in
#    OPENWRT_BUILDER_*_DIR)
#      printf "FOR [%s] CREATE [%s]\n" "$k" "$v"
#      [ -n "$v" ] || continue
#      mkdir -p -- "$v"
#    ;;
#  esac
#done
mkdir -p "$OPENWRT_BUILDER_LISTS_DIR"

# copy examples
# mkdir -p "$EXAMPLES_DIR"
