# openwrt_builder/rootfs/etc/s6-overlay/s6-rc.d/openwrt-builder-init/up
#!/command/with-contenv sh
printf "\033[32mINFO\033[0m INIT step.\n" >&2
set -eu

# export env paths
set -a
. /etc/openwrt-builder/paths.env
set +a

env | grep '^OPENWRT_BUILDER_.*_DIR=' >&2

while IFS='=' read -r k v; do
  case "$k" in
    OPENWRT_BUILDER_*_DIR)
      printf "FOR [%s] CREATE [%s]\n" "$k" "$v" >&2
      [ -n "$v" ] || continue
      mkdir -p -- "$v"
    ;;
  esac
done <<'EOF'
'"$(env)"'
EOF

# copy examples
# mkdir -p "$EXAMPLES_DIR"
